# overlay/etc/profile.d/fix_tmp.sh
#!/bin/bash
chown -f -R $USER:$USER /tmp/.X11-unix
